# assets/systemd/fidget-spinner-ui.service.in
[Unit]
Description=Fidget Spinner central navigator

[Service]
Type=simple
WorkingDirectory=@UI_PATH@
Environment=HOME=@HOME@
Environment=PATH=@HOME@/.cargo/bin:@LOCAL_BIN_DIR@:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=@LOCAL_BIN_DIR@/fidget-spinner-cli ui serve --path @UI_PATH@ --bind @UI_BIND@
Restart=always
RestartSec=3

[Install]
WantedBy=default.target
